# root/etc/services.d/vpn-monitor/run
#!/usr/bin/with-contenv bash

echo "Starting REAL-TIME VPN kill switch monitor..."

# Function to kill qBittorrent immediately
kill_qbittorrent() {
    local reason="$1"
    echo "======================================"
    echo "KILL SWITCH ACTIVATED!"
    echo "Reason: $reason"
    echo "======================================"
    
    # Kill ALL qBittorrent processes immediately with SIGKILL (cannot be ignored)
    killall -9 qbittorrent-nox 2>/dev/null
    pkill -9 -f qbittorrent 2>/dev/null
    
    # Remove VPN ready marker so it can't restart
    rm -f /tmp/vpn_ready
    
    # Log the event
    echo "[$(date)] KILL SWITCH: $reason" >> /config/killswitch.log
}

# Start real-time interface monitor in background
# This uses 'ip monitor' which provides INSTANT notifications on interface changes
(
    echo "Starting real-time interface state monitor..."
    ip monitor link | while read -r line; do
        if echo "$line" | grep -q "wg0"; then
            if echo "$line" | grep -qi "DOWN\|DELETED\|state DOWN"; then
                kill_qbittorrent "wg0 interface went DOWN (detected via ip monitor)"
            fi
        fi
    done
) &

# Start real-time route monitor in background
# Detects if default route through wg0 disappears
(
    echo "Starting real-time route monitor..."
    ip monitor route | while read -r line; do
        # If we see a route deletion and wg0 is involved, check the interface
        if echo "$line" | grep -q "Deleted\|del"; then
            if ! ip link show wg0 up 2>/dev/null | grep -q "state UP"; then
                kill_qbittorrent "wg0 route deleted or interface down"
            fi
        fi
    done
) &

# Periodic comprehensive checks (backstop for edge cases)
while true; do
    sleep 5
    
    # Check 1: Is wg0 interface up?
    if ! ip link show wg0 2>/dev/null | grep -q "state UP"; then
        kill_qbittorrent "wg0 interface state is not UP"
        
        # Try to restart VPN
        echo "Attempting to restart WireGuard..."
        wg-quick down wg0 2>/dev/null
        sleep 2
        if wg-quick up wg0 2>/dev/null; then
            echo "WireGuard restarted, verifying connectivity..."
            sleep 3
            if timeout 10 curl -s --interface wg0 https://api.ipify.org >/dev/null 2>&1; then
                echo "VPN restored successfully"
                touch /tmp/vpn_ready
            fi
        fi
        continue
    fi
    
    # Check 2: Can we reach internet through wg0?
    if ! timeout 5 curl -s --interface wg0 https://1.1.1.1 >/dev/null 2>&1; then
        kill_qbittorrent "Cannot reach internet through wg0"
        
        # Mark VPN as not ready
        rm -f /tmp/vpn_ready
        
        # Try to fix
        wg-quick down wg0 2>/dev/null
        sleep 2
        wg-quick up wg0 2>/dev/null
        sleep 3
        
        if timeout 10 curl -s --interface wg0 https://1.1.1.1 >/dev/null 2>&1; then
            echo "VPN connectivity restored"
            touch /tmp/vpn_ready
        fi
        continue
    fi
    
    # Check 3: Is qBittorrent running without VPN ready?
    if pgrep -x qbittorrent-nox >/dev/null 2>&1; then
        if [ ! -f /tmp/vpn_ready ]; then
            kill_qbittorrent "qBittorrent running without VPN ready marker"
        fi
    fi
    
    # Check 4: Verify WireGuard handshake is recent (within last 3 minutes)
    if command -v wg >/dev/null; then
        LAST_HANDSHAKE=$(wg show wg0 latest-handshakes 2>/dev/null | awk '{print $2}')
        if [ -n "$LAST_HANDSHAKE" ]; then
            CURRENT_TIME=$(date +%s)
            TIME_DIFF=$((CURRENT_TIME - LAST_HANDSHAKE))
            
            # If no handshake in 180 seconds (3 minutes), something is wrong
            if [ $TIME_DIFF -gt 180 ]; then
                kill_qbittorrent "WireGuard handshake stale (${TIME_DIFF}s old)"
            fi
        fi
    fi
done
